# deploy/systemd/gpu-usage-ntfy.service
[Unit]
Description=GPU idle ntfy monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__ROOT_DIR__
EnvironmentFile=-__ENV_PATH__
Environment=RUST_LOG=info
ExecStart=__BINARY_PATH__ --config __CONFIG_PATH__
Restart=always
RestartSec=5
UMask=0077
NoNewPrivileges=true
RestrictSUIDSGID=true
PrivateTmp=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
